# containers/quadlets/pf-runner.container
# pf-runner Container Quadlet
# The main pf task runner service
[Unit]
Description=pf Task Runner
Documentation=https://github.com/P4X-ng/pf-web-poly-compile-helper-runner
After=pf-web.network

[Container]
ContainerName=pf-runner
Image=localhost/pf-runner:latest
Pod=pf-web-api.pod

# Mount source code
Volume=pf-project-src.volume:/app/src:ro
Volume=pf-wasm-output.volume:/app/output:rw

# Environment configuration
Environment=PFY_FILE=/app/Pfyfile.pf
Environment=TERM=xterm-256color

# Security settings
SecurityLabelDisable=false
NoNewPrivileges=true

# Resource limits
MemoryMax=1G

# Keep container running for task execution
Exec=sleep infinity

[Service]
Restart=on-failure
TimeoutStartSec=90

[Install]
WantedBy=default.target
