# quadlet/pf-security-service.container
[Unit]
Description=pf Security Tools Service Container
Documentation=https://github.com/containers/podman/blob/main/docs/source/markdown/podman-systemd.unit.5.md
Requires=pf-main-pod.service
After=pf-main-pod.service
BindsTo=pf-main-pod.service

[Container]
# Container name
ContainerName=pf-security-service

# Image to use
Image=localhost/pf-security-tools:latest

# Pod assignment
Pod=pf-main-pod

# Volume mounts
Volume=pf-workspace:/workspace:rw,z
Volume=pf-builds:/builds:rw,z
Volume=pf-cache:/cache:rw,z

# Environment variables
Environment=WORKSPACE=/workspace

# Working directory
WorkingDir=/workspace

# User
User=1000:1000

# Security options - More permissive for debugging tools
SecurityLabelDisable=false
NoNewPrivileges=false

# Additional capabilities for debugging
AddCapability=SYS_PTRACE
AddCapability=SYS_ADMIN

# Resource limits
Memory=1g
CPUQuota=100%

# Health check
HealthCmd=gdb --version && python3 -c "import pwn" && radare2 -v || exit 1
HealthInterval=60s
HealthTimeout=30s
HealthRetries=3
HealthStartPeriod=45s

# Labels
Label=app=pf-development
Label=component=security-service

# Restart policy
Restart=always

# Command (keep container running)
Exec=sleep infinity

[Install]
WantedBy=default.target